# rootfs-overlay/common/etc/init.d/dam-client
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command=/home/decode/go/bin/dam-client
pidfile=/var/run/dam-client.pid

description="Tor DAM client"

depend() {
	after dam-dir
}

start() {
	start-stop-daemon --start -u decode -g decode \
		--pidfile $pidfile -b -m -x $command
}
