# rootfs-overlay/arm/etc/cron.hourly/decode-update
#!/bin/sh
# stub for the squashfs update mechanism
