# rootfs-overlay/arm/etc/profile.d/zulujava.sh
export PATH="/usr/local/ezdk/bin:$PATH"
