# rootfs-overlay/common/etc/rc.local
#!/bin/sh
# rc.local for base images

[ -f /etc/ssh/ssh_host_rsa_key.pub ] || ssh-keygen -A

netdata

exit 0
